# automation/deploy/stec-crm-worker.service
[Unit]
Description=STEC CRM normalization worker
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/stec-crm-worker
Environment=PYTHONUNBUFFERED=1
Environment=STEC_WORKER_HOST=0.0.0.0
Environment=STEC_WORKER_PORT=8080
Environment=STEC_TEMPLATE_PATH=/opt/stec-crm-worker/data/Nomina inscritos CRM.xlsx
ExecStart=/opt/stec-crm-worker/automation/deploy/start_stec_worker.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
